# infra/ansible/roles/etcd/files/etcd.service
[Unit]
Description=etcd
After=network-online.target
Wants=network-online.target

[Service]
User=etcd
Group=etcd
ExecStart=/usr/local/bin/etcd

[Install]
WantedBy=multi-user.target